# systemd/media-dire-videos.automount
[Unit]
Description=Automount DIRE Videos Share

[Automount]
Where=/media/dire/videos
TimeoutIdleSec=30s

[Install]
WantedBy=multi-user.target

